# config/http.ini
[program:httpserver]
command=bash -c  'LD_LIBRARY_PATH=./lib ./http_server'
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
startsecs=5                   ; number of secs prog must stay running (def. 1)
startretries=3  
user=root                     ; max # of serial start failures (default 3)
directory=/home/centos/workspace/http_server
stdout_logfile=/home/centos/log/http_server.log
stderr_logfile=/home/centos/log/http_server_err.log



[program:lift_ctrl_board]
command=bash -c './lift_ctrl_board' 
autostart=true
autorestart=true              ; retstart at unexpected quit (default: true)
startsecs=5                   ; number of secs prog must stay running (def. 1)
startretries=3  
user=root                     ; max # of serial start failures (default 3)
directory=/home/centos/workspace/http_server/test
stdout_logfile=/home/centos/log/http_test.log
stderr_logfile=/home/centos/log/http_test_err.log

